# config/files/system/etc/profile.d/more-modern-aliases.sh
alias .h='cd $(zoxide query --interactive)'
alias cat='bat --style plain'
alias ls='lsd'
alias lt='lsd -l --tree'
